# examples/systemd/minecraft_exporter.socket
[Unit]
Description=Minecraft Exporter

[Socket]
ListenStream=9940

[Install]
WantedBy=sockets.target
